# screenlogger.crontab
*/2 * * * * DISPLAY=:0.0 /path/to/screenlog.py
